# authman-samba-app/src/assembly/debian9/control/preinst
#!/bin/bash

if [ -e /etc/systemd/system/smb-con.service ]; then
    /bin/systemctl stop smb-con.service
    /bin/systemctl disable smb-con.service
fi

exit 0