# authman-samba-app/src/assembly/debian9/control/prerm
#!/bin/bash

/bin/systemctl stop bec.service
/bin/systemctl disable bec.service

exit 0